# root.cron
0,30 * * * *	docker system prune --all --force
